# php/supervisor/conf.d/php-xdebug.ini
[program:php-xdebug]
command=/usr/local/sbin/php-fpm -y /usr/local/etc/php-fpm.d/fpm-xdebug.conf -c /usr/local/etc/php/php-xdebug.ini -F
redirect_stderr=true

stdout_logfile=/dev/fd/2
stdout_logfile_maxbytes=0

stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0